# scripts/launchd/com.gurubet.evict-ghost-dry.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gurubet.evict-ghost-dry</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/wehandle/Projetos/pessoal/bets-estatistica/scripts/run-evict-ghost.sh</string>
        <string>--dry-run</string>
        <string>--notify-admin</string>
    </array>

    <!-- 2026-05-01 17:00 UTC = 14:00 BRT (one hour before the apply window). -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Month</key>
        <integer>5</integer>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>17</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/evict-ghost-dry.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/evict-ghost-dry.err</string>

    <!-- RunAtLoad=false: we only want this on the scheduled date. -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
